# scripts/supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/research_army/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid

; ── Ollama (model server) ──────────────────────────────────────────────────
[program:ollama]
command=/usr/local/bin/ollama serve
user=ubuntu
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/research_army/ollama.log
stderr_logfile=/var/log/research_army/ollama.err
environment=HOME="/home/ubuntu",OLLAMA_HOST="0.0.0.0:11434",OLLAMA_MAX_LOADED_MODELS="1",OLLAMA_NUM_PARALLEL="1"

; ── Redis ──────────────────────────────────────────────────────────────────
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=redis
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/research_army/redis.log
stderr_logfile=/var/log/research_army/redis.err

; ── Weaviate ───────────────────────────────────────────────────────────────
[program:weaviate]
command=/usr/bin/docker start -a weaviate
user=ubuntu
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/var/log/research_army/weaviate.log
stderr_logfile=/var/log/research_army/weaviate.err

; ── Research Army API ──────────────────────────────────────────────────────
[program:research_army]
command=/home/ubuntu/venv/bin/python main.py
directory=/home/ubuntu/research_army
user=ubuntu
autostart=true
autorestart=true
startsecs=15
startretries=5
stdout_logfile=/var/log/research_army/api.log
stderr_logfile=/var/log/research_army/api.err
environment=
    HOME="/home/ubuntu",
    PATH="/home/ubuntu/venv/bin:/usr/local/cuda/bin:/usr/local/bin:/usr/bin:/bin",
    VIRTUAL_ENV="/home/ubuntu/venv",
    PYTHONUNBUFFERED="1",
    CUDA_VISIBLE_DEVICES="0"

; Start order: redis and weaviate before api
[group:research_army]
programs=ollama,redis,weaviate,research_army
priority=999
